# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/storage-volume-bind-mount/templates/bind-mount.service.j2
[Unit]
DefaultDependencies=no
Description=Mount {{ src }} to {{ dest }}
Conflicts=umount.target
Before=local-fs.target umount.target

[Mount]
What={{ src }}
Where={{ dest }}
Type=none
Options=bind

[Install]
WantedBy=local-fs.target
