# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/coreos-mount/templates/mount.service.j2
[Unit]
Description=Mount {{ block_dev }} on {{ brick_mount_path }}

[Mount]
What={{ block_dev }}
Where={{ brick_mount_path }}
Type=xfs
Options=noatime

[Install]
WantedBy = multi-user.target
