# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/coreos-mount/templates/mount.service.j2
[Unit]
Description=Mount {{ block_dev }} on {{ mount_path }}
After=local-fs.target
[Mount]
What={{ block_dev }}
Where={{ mount_path }}
Type=xfs
Options=noatime
