# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/post-k8s-fixes/templates/mount.service.j2
[Unit]
Description=Mount {{ block_dev }} on {{ brick_mount_path }}
After=local-fs.target
[Mount]
What={{ block_dev }}
Where={{ brick_mount_path }}
Type=xfs
Options=noatime
