# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/openstack-storage-volume/templates/mount.service.j2
[Unit]
Description=Mount {{ item.src }} on {{ item.dest }}
After=local-fs.target

[Mount]
What={{ item.src }}
Where={{ item.dest }}
Type=xfs
Options=noatime

[Install]
WantedBy=multi-user.target
