# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/openstack-storage-volume/templates/mount.service.j2
[Unit]
Description=Mount {{ src }} on {{ dest }}
After=local-fs.target

[Mount]
What={{ src }}
Where={{ dest }}
Type=xfs
Options=noatime

[Install]
WantedBy=multi-user.target
