# FILES.deploy-tools/usr/local/lib/ndslabs/ansible/roles/openstack-storage-volume/templates/xfs-mount.service.j2
[Unit]
Description=Mount {{ storage_volume_dev }} on /media/storage
After=local-fs.target
 
 [Mount]
What={{ storage_volume_dev }}
Where=/media/storage
Type=xfs
Options=noatime
 
 [Install]
WantedBy=multi-user.target
